# automation/com.enriquezOS.n8n-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enriquezOS.n8n-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/allenenriquez/Developer/Allen-Enriquez/automation/run_n8n_monitor.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/n8n-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/n8n-monitor-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
