# automation/com.enriquezOS.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enriquezOS.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/allenenriquez/Desktop/Allen\ Enriquez/tools/dashboard &amp;&amp; exec /Library/Developer/CommandLineTools/usr/bin/python3 -m gunicorn -c gunicorn.conf.py app:app</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/enriquez-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/enriquez-dashboard-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/allenenriquez/Desktop/Allen Enriquez</string>
</dict>
</plist>
