# automation/com.enriquezOS.ryan-self-improve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enriquezOS.ryan-self-improve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/allenenriquez/Developer/Allen-Enriquez/automation/run_ryan_self_improve.sh</string>
    </array>

    <!-- Sunday 08:00 AM local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/ryan-self-improve.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/ryan-self-improve-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
