# automation/com.enriquezOS.automation-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enriquezOS.automation-status</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/allenenriquez/Developer/Allen-Enriquez/tools/automation_status.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/automation-status.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/automation-status-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez</string>
</dict>
</plist>
