# automation/com.enriquezOS.eod-ops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enriquezOS.eod-ops</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/usr/bin/python3 tools/eps/eod_ops_manager.py &amp;&amp; /usr/bin/python3 tools/eps/crm_sync.py &amp;&amp; /usr/bin/python3 tools/eps/check_outcomes.py</string>
    </array>

    <!-- Daily at 5:00 PM AEST -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/eod-ops.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez/.tmp/eod-ops-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/allenenriquez/Developer/Allen-Enriquez</string>
</dict>
</plist>
